# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:fastapi]
command=/app/.venv/bin/uvicorn api.main:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
stdout_logfile=/tmp/fastapi.out.log
stderr_logfile=/tmp/fastapi.err.log

[program:nextjs]
command=node_modules/.bin/next start -p 7860
directory=/app/frontend
autostart=true
autorestart=true
stdout_logfile=/tmp/nextjs.out.log
stderr_logfile=/tmp/nextjs.err.log
